# irrigationd.service
[Unit]
Description=Irrigation Controller Daemon
After=network.target

[Service]
# Path to the compiled binary
ExecStart=/usr/local/bin/irrigationd

# Required environment variables
Environment=IRRIGATIOND_TOKEN=changeme-very-secret-token
Environment=IRRIGATIOND_BIND_ADDR=0.0.0.0
Environment=IRRIGATIOND_LOG_LEVEL=debug

# Run as irrigationd user (remove if you want to test as root first)
User=irrigationd
Group=irrigationd

# Restart if it crashes
Restart=always
RestartSec=2

KillMode=process
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target

