# systemd/irrigationd.service
[Unit]
Description=Irrigation Controller Daemon
After=network.target

[Service]
# Path to binary
ExecStart=/usr/local/bin/irrigationd

# Run as dedicated unprivileged user/group
User=irrigationd
Group=irrigationd

# Load environment file with IRRIGATIOND_TOKEN
EnvironmentFile=-/etc/default/irrigationd

# Restart policy
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
SystemCallFilter=@default @network-io @file-system
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target